# test/fixtures/upstart/upstart-test-web-1.conf
#upstart-test-web-1
start on starting upstart-test-web
stop  on stopping upstart-test-web
respawn

env port="80"
env bind="0.0.0.0"
env mysql_host="localhost"
env mysql_user="toast"
env mysql_pass="peanutbutter"
env PATH="TEST_PATH"
env PORT="5000"
env FOREMAN_WORKER_NAME="web.1"


chdir TEST_DIR
setuid test-user
setgid test-user

exec node server.js 5000 >> /var/log/upstart-test/web-1.log 2>&1
