# test/fixtures/upstart/upstart-test-api-1.conf
#upstart-test-api-1
start on starting upstart-test-api
stop  on stopping upstart-test-api
respawn

env port="80"
env bind="0.0.0.0"
env mysql_host="localhost"
env mysql_user="toast"
env mysql_pass="peanutbutter"
env PATH="TEST_PATH"
env PORT="5100"
env FOREMAN_WORKER_NAME="api.1"


chdir TEST_DIR
setuid test-user
setgid test-user

exec node server.js 10000 >> /var/log/upstart-test/api-1.log 2>&1
