# test/fixtures/upstart/upstart-test-log.conf
#upstart-test-log
start on starting upstart-test
stop on stopping upstart-test

