# test/fixtures/upstart/upstart-test-web.conf
#upstart-test-web
start on starting upstart-test
stop on stopping upstart-test

