# test/fixtures/systemd/systemd-test-log.target
[Unit]
Requires=systemd-test.target
Wants=systemd-test-log-.service
