# test/fixtures/systemd/systemd-test.target
[Unit]
Wants= systemd-test-web.target systemd-test-api.target systemd-test-log.target

[Install]
WantedBy=multi-user.target
