# test/fixtures/systemd/systemd-test-web.target
[Unit]
Requires=systemd-test.target
Wants= systemd-test-web-1.service
