# lib/systemd/foreman-APP.target
[Unit]
Requires={{{application}}}.target
Wants={{#numbers}} {{{application}}}-{{{process}}}-{{{number}}}.service{{/numbers}}
