# lib/systemd/foreman.target
[Unit]
{{{#processes}}}Wants={{{application}}}-{{{process}}}.target
{{{/processes}}}