# lib/systemd/foreman.target
[Unit]
Wants={{#processes}} {{{application}}}-{{{process}}}.target{{/processes}}

[Install]
WantedBy=multi-user.target
