# lib/systemd/foreman-APP.target
[Unit]
Requires={{{application}}}.target
{{{#numbers}}}Wants={{{application}}}-{{{process}}}-{{{number}}}.service
{{{/numbers}}}