# lib/upstart/foreman-APP-N.conf
#{{{application}}}-{{{process}}}-{{{number}}}
start on starting {{{application}}}-{{{process}}}
stop  on stopping {{{application}}}-{{{process}}}
respawn

{{#envs}}env {{{key}}}="{{{value}}}"
{{/envs}}

chdir {{{cwd}}}
setuid {{{user}}}
setgid {{{group}}}

exec {{{command}}} >> {{{logs}}}/{{{application}}}/{{{process}}}-{{{number}}}.log 2>&1
