# lib/upstart/foreman-APP.conf
#{{{application}}}-{{{process}}}
start on starting {{{application}}}
stop on stopping {{{application}}}

