# lib/upstart/foreman.conf
#{{{application}}}
pre-start script

bash << "EOF"
  mkdir -p {{{logs}}}/{{{application}}}
  chown -R {{{user}}} {{{logs}}}/{{{application}}}
EOF

end script

start on runlevel [2345]

stop on runlevel [016]
