# test/fixtures/procfile-comment/foreman-run-1.conf
#foreman-run-1
start on starting foreman-run
stop  on stopping foreman-run
respawn

env PATH="/usr/local/lib/node_modules/npm/bin/node-gyp-bin:/home/blittle/dev/node-foreman/node_modules/.bin:TEST_PATH:/usr/lib64/mpi/gcc/openmpi/bin:/home/blittle/bin:TEST_PATH:/usr/bin:/bin:/usr/bin/X11:/usr/X11R6/bin:/usr/games:/home/blittle/bin/gradle/bin:/home/blittle/bin/byzanz:/home/blittle/bin/apache-maven/bin:/home/blittle/bin/mongodb/bin"
env PORT="5000"
env FOREMAN_WORKER_NAME="run.1"


chdir TEST_DIR
setuid root
setgid root

exec touch should-exist.js >> /var/log/foreman/run-1.log 2>&1
