# test/fixtures/procfile-comment/foreman-run.conf
#foreman-run

# Generated by node-foreman

description "foreman-run service"

# Automatic start/stop based on foreman, which starts on bootup
# To manually start/stop only this job:
# /sbin/initctl start foreman-run
# /sbin/initctl stop foreman-run
start on starting foreman
stop  on stopping foreman
respawn

# Configuration for this daemon
env PATH="something-totally-bogus"


# Not supported by older versions of Upstart, like on RHEL/CentOS
chdir TEST_DIR
setuid root
setgid root

# Arbitrarily high limit on number of open file descriptors
limit nofile 50000 50000

# Enable core dumps
limit core unlimited unlimited

# How to run the application
script

    # Abort if some part of this script fails
    set -e

    # Temporary fifo for sending log output to
    mkfifo /tmp/foreman-run-log-fifo

    # A logger instance that pipes the fifo to syslog
    ( logger -t foreman-run </tmp/foreman-run-log-fifo & )

    # Redirect the current process's stdout to the logger
    exec >/tmp/foreman-run-log-fifo

    # Clean up the fifo so we don't leave it laying around
    rm /tmp/foreman-run-log-fifo

    # Run our actual application
    #  - exec used instead of setuid, see comment above
    #  - stderr redirection is done last so that we don't lose error messages
    #    related to the Upstart config itself
    # Uncomment and use this line if you are using an older version of Upstart
    #exec su -m -s /bin/sh -c 'exec "$0" "$@"' root -- npm start 2>&1
    exec touch should-exist.js 2>&1

end script
