# lib/supervisord/foreman-APP-N.conf
[program:{{{application}}}-{{{process}}}-{{{number}}}]
command={{{ command }}}
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile={{{logs}}}/{{{application}}}-{{{process}}}-{{{number}}}.stdout.log
stderr_logfile={{{logs}}}/{{{application}}}-{{{process}}}-{{{number}}}.stderr.log
user={{{user}}}
directory={{{cwd}}}
environment={{{envs}}}

