# test/fixtures/supervisord/supervisord-test-api-1.conf
[program:supervisord-test-api-1]
command=node server.js 10000
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/supervisord-test-api-1.stdout.log
stderr_logfile=/var/log/supervisord-test-api-1.stderr.log
user=test-user
directory=TEST_DIR
environment=port="80",bind="0.0.0.0",mysql_host="localhost",mysql_user="toast",mysql_pass="peanutbutter",PATH="TEST_PATH",PORT="5100",FOREMAN_WORKER_NAME="api.1"

