# test/fixtures/supervisord/supervisord-test-web-1.conf
[program:supervisord-test-web-1]
command=node server.js 5000
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/supervisord-test-web-1.stdout.log
stderr_logfile=/var/log/supervisord-test-web-1.stderr.log
user=test-user
directory=TEST_DIR
environment=port="80",bind="0.0.0.0",mysql_host="localhost",mysql_user="toast",mysql_pass="peanutbutter",PATH="TEST_PATH",PORT="5000",FOREMAN_WORKER_NAME="web.1"

