# lib/upstart-single/foreman-APP.conf
#{{{application}}}-{{{process}}}
start on starting {{{application}}}-{{{process}}}
stop  on stopping {{{application}}}-{{{process}}}
respawn
{{#envs}}
env {{{key}}}={{{value}}}{{/envs}}

export PATH

chdir {{{cwd}}}
setuid {{{user}}}
setgid {{{group}}}
limit nofile 50000 50000

exec {{{command}}} >> {{{logs}}}/{{{application}}}/{{{process}}}.log 2>&1
