# test/fixtures/upstart-single/upstart-single-test.conf
#upstart-single-test
pre-start script

bash << "EOF"
  mkdir -p /var/log/upstart-single-test
  chown -R root /var/log/upstart-single-test
EOF

end script

start on (started network-interface
          or started network-manager
          or started networking)

stop on (stopping network-interface
         or stopping network-manager
         or stopping networking)

