# test/fixtures/upstart-single/upstart-single-test.conf
#upstart-single-test

description "upstart-single-test super launcher"

start on runlevel [2345]

stop on runlevel [016]
