# test/fixtures/upstart-single/upstart-single-test.conf
#upstart-single-test

description "upstart-single-test super launcher"

# Start on boot/reboot up
start on runlevel [2345]

# Stop on shutdown/reboot
stop on runlevel [016]

# To start/stop all workers under this app manually you can run:
#
# /sbin/initctl start upstart-single-test
# /sbin/initctl stop upstart-single-test
#
# For instructions on starting/stopping the individual components of this app,
# see the start/stop descriptions in their Upstart job configurations.
