# test/fixtures/upstart-single/upstart-single-test-log.conf
#upstart-single-test-log
start on starting upstart-single-test-log
stop  on stopping upstart-single-test-log
respawn

env port=80
env bind=0.0.0.0
env mysql_host=localhost
env mysql_user=toast
env mysql_pass=peanutbutter
env PATH=/usr/local/bin
env PORT=5200

export PATH

chdir /Users/ryan/work/node-foreman/test
setuid root
setgid root
limit nofile 50000 50000

exec node server.js 100000 >> /var/log/upstart-single-test/log.log 2>&1
