# test/fixtures/smf/smf-test-log.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="smf-test">
  <service name="site/smf-test" type="service" version="1">

    <create_default_instance enabled="true"/>

    <single_instance/>

    <dependency name="network" grouping="require_all" restart_on="refresh" type="service">
      <service_fmri value="svc:/milestone/network:default"/>
    </dependency>

    <dependency name="filesystem" grouping="require_all" restart_on="refresh" type="service">
      <service_fmri value="svc:/system/filesystem/local"/>
    </dependency>

    <method_context working_directory="TEST_DIR">
      <method_credential user="test-user" group="test-user" privileges="basic,net_privaddr"/>
      <method_environment>
        <envvar name="PATH" value="TEST_PATH"/>
<envvar name="bind" value="0.0.0.0"/>
<envvar name="mysql_host" value="localhost"/>
<envvar name="mysql_pass" value="peanutbutter"/>
<envvar name="mysql_user" value="toast"/>
<envvar name="port" value="80"/>
      </method_environment>
    </method_context>

    <exec_method
      type="method"
      name="start"
      exec="node server.js 100000"
      timeout_seconds="60"/>

    <exec_method
      type="method"
      name="stop"
      exec=":kill"
      timeout_seconds="60"/>

    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="child"/>
      <propval name="ignore_error" type="astring" value="core,signal"/>
    </property_group>

    <property_group name="application" type="application">

    </property_group>


    <stability value="Evolving"/>

    <template>
      <common_name>
        <loctext xml:lang="C">smf-test</loctext>
      </common_name>
    </template>

  </service>

</service_bundle>
